# deploy/aura.service
[Unit]
Description=Aura Task Manager
After=network.target
Wants=network.target

[Service]
Type=simple
User=clevercode
Group=clevercode
WorkingDirectory=/home/clevercode/aura/backend
ExecStart=/home/clevercode/aura/backend/bin/aura
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=aura

EnvironmentFile=-/etc/aura/env
Environment=PORT=9001
Environment=DB_PATH=/home/clevercode/aura/backend/data/aura.db
Environment=FRONTEND_DIR=/home/clevercode/aura/frontend/build
Environment=APP_URL=https://blackbox.clevercode.ts.net
Environment=FRONTEND_URL=https://blackbox.clevercode.ts.net
Environment=ENV=production

[Install]
WantedBy=multi-user.target
